# roles/ecs-agent/templates/ecs-agent.path.j2
[Unit]
Description=Manage path dependency for ecs-agent


[Path]
DirectoryNotEmpty=/opt/nginx_conf
