# roles/ecs-agent/templates/ecs-agent.path.j2
[Path]
DirectoryNotEmpty=/opt/nginx_conf
