# roles/ecs-agent/templates/ecs-agent.timer.j2
[Unit]
Description=AWS ECS Agent delay timer

[Timer]
OnActiveSec=300

[Install]
WantedBy=multi-user.target
